# config/systemd/user/check-battery.service
[Unit]
Description=check battery

[Service]
Type=oneshot
ExecStart=/home/picuber/.config/systemd/user/check-battery.sh
